# deploy/supervisor-reverb.conf
# ============================================
# Supervisor — Laravel Reverb WebSocket
# Letakkan di: /etc/supervisor/conf.d/absensi-reverb.conf
# ============================================

[program:absensi-reverb]
process_name=%(program_name)s
command=php /var/www/absensi/artisan reverb:start --port=8080
autostart=true
autorestart=true
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/absensi/storage/logs/reverb.log
stopwaitsecs=10
stopsignal=TERM
