# deploy/supervisor-reverb.conf
# ============================================
# Supervisor — Laravel Reverb WebSocket
# Letakkan di: /etc/supervisor/conf.d/sistem-kehadiran-reverb.conf
# ============================================

[program:sistem-kehadiran-reverb]
process_name=%(program_name)s
command=php /var/www/sistem-kehadiran/artisan reverb:start --port=8080
autostart=true
autorestart=true
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/sistem-kehadiran/storage/logs/reverb.log
stopwaitsecs=10
stopsignal=TERM
